# web_manager.service
[Unit]
Description=Attendance Web Manager UI
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/ai
Environment="PATH=/home/pi/ai/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="VIRTUAL_ENV=/home/pi/ai/venv"
ExecStart=/home/pi/ai/venv/bin/python /home/pi/ai/web_manager.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
